# docker/rootfs/etc/s6-overlay/s6-rc.d/hbbr/run
#!/command/with-contenv sh
cd /data
/usr/bin/hbbr
