# docker/rootfs/etc/s6-overlay/s6-rc.d/hbbr/run
#!/command/execlineb -P
posix-cd /data
/usr/bin/hbbr
