# docker/rootfs/etc/s6-overlay/s6-rc.d/hbbr/run
#!/command/with-contenv sh
cd /data || return
PARAMS=
[ "${ENCRYPTED_ONLY}" = "1" ] && PARAMS="-k _"
if [ "$(id -u)" -ne 0 ] ; then
  /usr/bin/hbbr $PARAMS
else
  exec s6-setuidgid rustdesk /usr/bin/hbbr $PARAMS
fi
