# docker/rootfs/etc/s6-overlay/s6-rc.d/key-secret/up
/etc/s6-overlay/s6-rc.d/key-secret/up.real
