# docker/rootfs/etc/s6-overlay/s6-rc.d/hbbs/run
#!/command/execlineb -P
posix-cd /data
/usr/bin/hbbs -r $RELAY
