# docker/rootfs/etc/s6-overlay/s6-rc.d/hbbs/run
#!/command/with-contenv sh
sleep 2
cd /data || return
PARAMS=
[ "${ENCRYPTED_ONLY}" = "1" ] && PARAMS="-k _"
exec s6-setuidgid rustdesk /usr/bin/hbbs -r $RELAY $PARAMS
